# deploy/launchd/com.aml.nova.agent-core.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aml.nova.agent-core</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/node</string>
    <string>__NOVA_PATH__/apps/agent-core/dist/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__NOVA_PATH__</string>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__NOVA_PATH__/data/logs/agent-core.out.log</string>
  <key>StandardErrorPath</key>
  <string>__NOVA_PATH__/data/logs/agent-core.err.log</string>
</dict>
</plist>
